# .config/systemd/user/check-battery.service
[Unit]
Description=Check battery

[Service]
ExecStart=/home/ali/.local/bin/miscellaneous/check-battery
