# .config/systemd/user/check-battery.service
[Unit]
Description=Check battery

[Service]
ExecStart=/home/ali/g/scripts/miscellaneous/check-battery
