# .config/systemd/user/check-battery.service
[Unit]
Description=Check battery

[Service]
ExecStart=check-battery
