# .config/profile.d/20-rclone
#!/bin/sh

export RCLONE_PASSWORD_COMMAND="rbw get rclone"
