# .config/profile.d/20-rclone
#!/bin/sh

export RCLONE_PASSWORD_COMMAND="secret-tool lookup program rclone user $USER"
