# systemd-units/check-battery.service
[Unit]
Description=Check battery

[Service]
ExecStart=/home/ali/g/scripts/miscellaneous/battery.sh
