# udev-rules/98-discharging.rules
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", RUN+="/usr/bin/sh /home/ali/1/programs/bash/bateria2"
